# deploy/longzhongxian.service
[Unit]
Description=笼中仙 AI 招聘助手
After=network.target postgresql.service

[Service]
Type=simple
User=ycm
WorkingDirectory=/home/ycm/longzhongxian
EnvironmentFile=/home/ycm/longzhongxian/.env.local
ExecStart=/home/ycm/longzhongxian/.venv/bin/uvicorn src.app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
